# F5L049v1000_1.00.04_GPL/F5L049v1000_1.00.04_OSS/preinstall/etc/inittab
::sysinit:/etc/rc.d/rc.S
::restart:/etc/rc.d/rc.6
::shutdown:/etc/rc.d/rc.0
::shutdown:/bin/umount -r /
